# files/xmonad/xmonad.target
[Unit]
Description=Xmonad running
BindsTo=graphical-session.target
